# src/systemd/mysqlarabicci.service
[Unit]
Description=Inject Arabic utf8 collation into MySQL charset index

[Service]
Type=oneshot
ExecStart=/usr/local/bin/mysqlarabicci --inject
